# packaging/DEBIAN/postinst
#!/bin/bash
set -e

CONFIG_DIR=/etc/patchpanel
DATA_DIR=/var/lib/patchpanel
LOG_DIR=/var/log/patchpanel
INSTALL_DIR=/opt/patchpanel
TEMPLATE_DIR=${INSTALL_DIR}/config-templates

# 1. System user + group
if ! getent passwd patchpanel > /dev/null 2>&1; then
    adduser --system --group --home "${INSTALL_DIR}" \
            --no-create-home --shell /usr/sbin/nologin \
            --gecos "PatchPanel system user" patchpanel
fi

# 2. Supplementary group: haproxy (so we can rewrite /etc/haproxy/haproxy.cfg)
if getent group haproxy > /dev/null 2>&1; then
    usermod -aG haproxy patchpanel || true
fi

# 2b. Keepalived integration. If keepalived is installed on this host
#     ensure patchpanel can render to /etc/keepalived/keepalived.conf.
#     Most distros don't ship a `keepalived` group, so we set the
#     directory mode + ownership directly. SIGUSR2-based VRRP state
#     detection also needs CAP_KILL — granted via systemd unit
#     AmbientCapabilities, not here.
if [[ -d /etc/keepalived ]]; then
    if getent group keepalived > /dev/null 2>&1; then
        usermod -aG keepalived patchpanel || true
    fi
    chgrp patchpanel /etc/keepalived || true
    chmod g+rwx /etc/keepalived || true
    if [[ -f /etc/keepalived/keepalived.conf ]]; then
        chgrp patchpanel /etc/keepalived/keepalived.conf || true
        chmod g+rw /etc/keepalived/keepalived.conf || true
    fi
fi

# 3. Directories
install -d -o patchpanel -g patchpanel -m 755 "${DATA_DIR}"
install -d -o patchpanel -g patchpanel -m 755 "${LOG_DIR}"
install -d -o patchpanel -g patchpanel -m 755 "${CONFIG_DIR}"
install -d -o patchpanel -g patchpanel -m 700 "${CONFIG_DIR}/ssl"

# 4. Fresh install vs upgrade
case "$1" in
    configure)
        # 4a. Copy config template only if user config doesn't exist
        if [[ ! -f "${CONFIG_DIR}/config.yaml" ]]; then
            cp "${TEMPLATE_DIR}/production-config.yaml" "${CONFIG_DIR}/config.yaml"
            chown patchpanel:patchpanel "${CONFIG_DIR}/config.yaml"
            chmod 640 "${CONFIG_DIR}/config.yaml"
            FRESH_INSTALL=1
        else
            FRESH_INSTALL=0
        fi

        # 4b. JWT secret — generate if missing, substitute placeholder in config
        if [[ ! -f "${CONFIG_DIR}/.jwt-secret" ]]; then
            openssl rand -hex 32 > "${CONFIG_DIR}/.jwt-secret"
            chown patchpanel:patchpanel "${CONFIG_DIR}/.jwt-secret"
            chmod 600 "${CONFIG_DIR}/.jwt-secret"
        fi

        # Substitute placeholder in config (only matters on fresh install;
        # on upgrade the user's config already has the real secret)
        if grep -q '__JWT_SECRET_FROM_FILE__' "${CONFIG_DIR}/config.yaml" 2>/dev/null; then
            JWT_SECRET=$(cat "${CONFIG_DIR}/.jwt-secret")
            sed -i "s|__JWT_SECRET_FROM_FILE__|${JWT_SECRET}|g" "${CONFIG_DIR}/config.yaml"
        fi

        # 4c. Setup token for first-run wizard (fresh install only)
        if [[ "${FRESH_INSTALL}" = "1" ]] && [[ ! -f "${CONFIG_DIR}/setup.token" ]]; then
            openssl rand -hex 32 > "${CONFIG_DIR}/setup.token"
            chown patchpanel:patchpanel "${CONFIG_DIR}/setup.token"
            chmod 600 "${CONFIG_DIR}/setup.token"
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        # No action needed for these dpkg phases.
        ;;
    *)
        echo "postinst called with unknown argument: $1" >&2
        exit 1
        ;;
esac

# 5. Reload systemd — skip when there's no running systemd (Docker builds,
# chroots, schroot, etc.). Lintian recommends this pattern; otherwise
# `dpkg --configure patchpanel` fails inside any non-systemd container.
if [[ -d /run/systemd/system ]]; then
    systemctl daemon-reload
fi

# 6. Banner
SETUP_TOKEN=""
if [[ "${FRESH_INSTALL:-0}" = "1" ]] && [[ -f "${CONFIG_DIR}/setup.token" ]]; then
    SETUP_TOKEN=$(cat "${CONFIG_DIR}/setup.token")
fi

HOSTNAME_FQDN=$(hostname -f 2>/dev/null || hostname)

echo ""
echo "================================================================="
echo "  PatchPanel installed."
echo ""
echo "  Config:    ${CONFIG_DIR}/config.yaml"
echo "  Data:      ${DATA_DIR}/"
echo "  Logs:      journalctl -fu patchpanel"
echo ""
echo "  Start:     systemctl enable --now patchpanel"
echo ""
if [[ -n "${SETUP_TOKEN}" ]]; then
    echo "  First-run setup URL (one-time, consumed by the wizard):"
    echo "    https://${HOSTNAME_FQDN}:8099/setup-admin?token=${SETUP_TOKEN}"
    echo ""
    echo "  Or if you forget the URL, the token alone is in:"
    echo "    ${CONFIG_DIR}/setup.token"
    echo ""
    echo "  Recovery (if locked out later):"
    echo "    patchpanel user-add --username admin"
    echo "    patchpanel user-reset --username admin"
fi
echo "================================================================="
echo ""

#DEBHELPER#
